# src-tauri/resources/ghostcat.service.template
[Unit]
Description=PLACEHOLDER_DESCRIPTION

[Service]
Type=oneshot
ExecStart=PLACEHOLDER_EXE_PATH
